# ubuntu-14.04-trusty-tahr/upstart.conf
# Mattermost: http://www.mattermost.org/
start on runlevel [2345]
stop on runlevel [016]
respawn
chdir /opt/mattermost
# TODO: fix
setuid mattermost
exec bin/platform
